# allconfig/etc/init.d/killpss
#!/sbin/openrc-run
# Copyright 2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="killpss daemon"
description="kill listed process"
supervisor="supervise-daemon"
healthcheck_delay=5
healthcheck_timer=1
respawn_max=0
command="/usr/bin/sleep"
command_args="infinity"

healthcheck() {
	for pls in ${PLIST[@]}; do
		local psid=$(ps -ef | tr -s " " | cut -d " " -f 2,8 | grep ${pls} | grep -v $0 | awk '{print $1}')
		[[ ! -n $psid ]] || return 1
	done
}

unhealthy() {
	for pls in ${PLIST[@]}; do
		local psid=$(ps -ef | tr -s " " | cut -d " " -f 2,8 | grep ${pls} | grep -v $0 | awk '{print $1}')
		[[ ! -n $psid ]] || /usr/bin/killall ${pls}
	done
}
