# allconfig/etc/local.d/brightness.start
# /etc/local.d/local.start:
# $Header: /home/cvsroot/gentoo-src/rc-scripts/etc/conf.d/local.start,v 1.4 2002/11/18 19:39:22

# This is a good place to load any misc.
# programs on startup ( 1>&2 )

# Adjust backlight brightness
#echo 30 > /sys/class/backlight/acpi_video0/brightness
xbacklight -ctrl acpi_video0 -set 40
#echo "boot at $(date), brightness:" >> /root/boot_brightness
#cat /sys/class/backlight/acpi_video0/brightness >> /root/boot_brightness
